# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

; ── HTTP server ──────────────────────────────────────────────────────────────
[program:php]
command=php -d variables_order=EGPCS /var/www/html/artisan serve --host=0.0.0.0 --port=80
user=sail
environment=LARAVEL_SAIL=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true

; ── Queue worker ─────────────────────────────────────────────────────────────
[program:queue-worker]
command=php /var/www/html/artisan queue:work --queue=ingest,resolution,high,default --sleep=3 --tries=3 --max-time=3600
user=sail
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
stopwaitsecs=3600

; ── Scheduler (every minute) ─────────────────────────────────────────────────
[program:scheduler]
command=bash -c "while true; do php /var/www/html/artisan schedule:run --no-interaction; sleep 60; done"
user=sail
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true

; ── Recorder SQLite sync (every 30 s) ────────────────────────────────────────
[program:ingest-sync]
command=bash -c "while true; do php /var/www/html/artisan ingest:sync --once; sleep 30; done"
user=sail
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
